# kakao-slack-bot.service
[Unit]
Description=KakaoTalk-Slack Profile Sync Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/kakao-slack-bot
ExecStart=/usr/bin/node src/index.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=kakao-slack-bot

# Environment variables
Environment=NODE_ENV=production

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
